# rootfs/etc/s6-overlay/s6-rc.d/apex-cloud-link/run
#!/usr/bin/with-contenv bashio
set -e

CONFIG_PATH=/data/options.json
FRP_VERSION="0.56.0"
FRP_SERVER_ADDR="cloud.apexinfosys.in"
FRP_SERVER_PORT=7000
FRP_DIR="/data/frp"
FRPC_BIN="${FRP_DIR}/frpc"
FRPC_CONFIG="/tmp/frpc.toml"
DEVICE_TOKEN_PATH="/data/cloud_device_token"
SSH_SESSION_STATE_PATH="/data/ssh_session_state.json"
SSH_SESSION_MARKERS_PATH="/data/ssh_session_markers"
SSH_SESSION_KEY_DIR="/data/ssh_sessions"

HEARTBEAT_INTERVAL=20
HEARTBEAT_PID=""
SSH_SYNC_INTERVAL=15
SSH_SYNC_PID=""
DEVICE_TOKEN=""

SUBDOMAIN=$(bashio::config 'subdomain')
ACCESS_TOKEN=$(bashio::config 'access_token')
FLEET_REPORTING_ENABLED=$(bashio::config 'fleet_reporting_enabled')
CLOUD_API_URL=$(bashio::config 'cloud_api_url')
DEVICE_UID_OPTION=$(bashio::config 'device_uid')
DEVICE_NAME_OPTION=$(bashio::config 'device_name')
SSH_REMOTE_PORT_OPTION=$(bashio::config 'ssh_remote_port')
SSH_LOCAL_PORT_OPTION=$(bashio::config 'ssh_local_port')
SSH_REMOTE_USER_OPTION=$(bashio::config 'ssh_remote_user')
SSH_AUTH_USER_OPTION=$(bashio::config 'ssh_auth_user')
SSH_AUTH_KEYS_FILE_OPTION=$(bashio::config 'ssh_authorized_keys_file')

SUBDOMAIN=${SUBDOMAIN//.cloud.apexinfosys.in/}

is_truthy() {
    case "$1" in
        1|true|TRUE|True|yes|YES|Yes|on|ON|On)
            return 0
            ;;
        *)
            return 1
            ;;
    esac
}

sanitize_uid() {
    printf '%s' "$1" | tr -cd 'a-zA-Z0-9._-'
}

is_valid_port() {
    case "$1" in
        ''|*[!0-9]*)
            return 1
            ;;
    esac

    if [ "$1" -lt 1 ] || [ "$1" -gt 65535 ]; then
        return 1
    fi

    return 0
}

sanitize_remote_user() {
    case "$1" in
        [a-zA-Z_][a-zA-Z0-9_-]*)
            printf '%s' "$1"
            ;;
        *)
            printf 'root'
            ;;
    esac
}

marker_name() {
    local session_id="$1"
    local fingerprint="$2"
    local username="$3"
    local safe_fingerprint
    safe_fingerprint=$(printf '%s' "$fingerprint" | tr -cd 'a-zA-Z0-9._:-')

    if [ -z "$safe_fingerprint" ]; then
        safe_fingerprint="unknown"
    fi

    printf 'session_%s__%s__%s' "$session_id" "$username" "$safe_fingerprint"
}

normalize_state_file() {
    if [ ! -f "$SSH_SESSION_STATE_PATH" ]; then
        printf '{"applied_session_ids":[],"active_markers":[]}' > "$SSH_SESSION_STATE_PATH"
    fi
}

read_applied_ids_csv() {
    normalize_state_file
    jq -r '.applied_session_ids // [] | map(tostring) | join(",")' "$SSH_SESSION_STATE_PATH"
}

write_state() {
    local applied_ids_csv="$1"
    local markers_csv="$2"

    jq -n \
        --arg applied_ids "$applied_ids_csv" \
        --arg markers "$markers_csv" \
        '{
            applied_session_ids: ($applied_ids | if . == "" then [] else split(",") | map(tonumber) end),
            active_markers: ($markers | if . == "" then [] else split(",") end)
        }' > "$SSH_SESSION_STATE_PATH"
}

ensure_ssh_storage() {
    mkdir -p "$SSH_SESSION_MARKERS_PATH"
    mkdir -p "$SSH_SESSION_KEY_DIR"
}

get_authorized_keys_path() {
    if [ -n "$SSH_AUTH_KEYS_FILE_OPTION" ] && [ "$SSH_AUTH_KEYS_FILE_OPTION" != "null" ]; then
        printf '%s' "$SSH_AUTH_KEYS_FILE_OPTION"
        return
    fi

    local user_name="$1"

    if [ "$user_name" = "root" ]; then
        printf '/root/.ssh/authorized_keys'
        return
    fi

    printf '/home/%s/.ssh/authorized_keys' "$user_name"
}

resolve_ssh_auth_user() {
    if [ -n "$SSH_AUTH_USER_OPTION" ] && [ "$SSH_AUTH_USER_OPTION" != "null" ]; then
        sanitize_remote_user "$SSH_AUTH_USER_OPTION"
        return
    fi

    sanitize_remote_user "$SSH_REMOTE_USER"
}

ensure_authorized_keys_parent() {
    local keys_path="$1"
    local parent_dir

    parent_dir=$(dirname "$keys_path")
    mkdir -p "$parent_dir"
    touch "$keys_path"
    chmod 700 "$parent_dir" || true
    chmod 600 "$keys_path" || true
}

get_local_ips_csv() {
    local ip_list=''

    if command -v ip >/dev/null 2>&1; then
        ip_list=$(ip -o -4 addr show scope global 2>/dev/null | awk '{split($4, parts, "/"); print parts[1]}' | paste -sd, -)
    fi

    if [ -z "$ip_list" ] && command -v hostname >/dev/null 2>&1; then
        ip_list=$(hostname -i 2>/dev/null | tr ' ' ',' | sed 's/,$//')
    fi

    printf '%s' "$ip_list"
}

cloud_post() {
    local endpoint="$1"
    local payload="$2"
    local token="$3"
    local base_url="${CLOUD_API_URL%/}"
    local url="${base_url}${endpoint}"

    if [ -z "$token" ]; then
        curl -fsS --max-time 12 \
            -H "Content-Type: application/json" \
            -X POST "$url" \
            -d "$payload" 2>/dev/null || true
        return
    fi

    curl -fsS --max-time 12 \
        -H "Content-Type: application/json" \
        -H "Authorization: Bearer ${token}" \
        -X POST "$url" \
        -d "$payload" 2>/dev/null || true
}

save_device_token() {
    if [ -z "$DEVICE_TOKEN" ]; then
        return
    fi

    printf '%s' "$DEVICE_TOKEN" > "$DEVICE_TOKEN_PATH"
    chmod 600 "$DEVICE_TOKEN_PATH" || true
}

send_device_log() {
    local level="$1"
    local event_type="$2"
    local message="$3"
    local payload_json="$4"

    if [ "$FLEET_ACTIVE" -ne 1 ] || [ -z "$DEVICE_TOKEN" ]; then
        return
    fi

    local payload
    if [ -n "$payload_json" ]; then
        payload=$(jq -n \
            --arg level "$level" \
            --arg event_type "$event_type" \
            --arg message "$message" \
            --argjson payload "$payload_json" \
            '{ level: $level, event_type: $event_type, message: $message, payload: $payload }')
    else
        payload=$(jq -n \
            --arg level "$level" \
            --arg event_type "$event_type" \
            --arg message "$message" \
            '{ level: $level, event_type: $event_type, message: $message }')
    fi

    cloud_post "/api/internal/devices/log" "$payload" "$DEVICE_TOKEN" >/dev/null
}

build_fleet_payload() {
    local state="$1"
    local local_ips_csv
    local local_ips_json='[]'
    local remote_port_json='null'

    local_ips_csv=$(get_local_ips_csv)

    if [ -n "$local_ips_csv" ]; then
        local_ips_json=$(jq -n --arg csv "$local_ips_csv" '$csv | split(",")')
    fi

    if [ "$SSH_REMOTE_PORT" -gt 0 ]; then
        remote_port_json="$SSH_REMOTE_PORT"
    fi

    jq -n \
        --arg access_token "$ACCESS_TOKEN" \
        --arg device_uid "$DEVICE_UID" \
        --arg device_name "$DEVICE_NAME" \
        --arg hostname "$DEVICE_HOSTNAME" \
        --argjson local_ips "$local_ips_json" \
        --argjson ssh_port "$SSH_LOCAL_PORT" \
        --arg remote_user "$SSH_REMOTE_USER" \
        --arg tunnel_host "$FRP_SERVER_ADDR" \
        --argjson tunnel_port "$remote_port_json" \
        --arg addon_version "apex-cloud-link" \
        --arg agent_state "$state" \
        '{
            access_token: $access_token,
            device_uid: $device_uid,
            device_name: $device_name,
            hostname: $hostname,
            local_ips: $local_ips,
            ssh_port: $ssh_port,
            remote_user: $remote_user,
            tunnel_host: $tunnel_host,
            tunnel_port: $tunnel_port,
            addon_version: $addon_version,
            agent_state: $agent_state
        }'
}

register_fleet_device() {
    if [ "$FLEET_ACTIVE" -ne 1 ]; then
        return
    fi

    local cached_token=''
    if [ -f "$DEVICE_TOKEN_PATH" ]; then
        cached_token=$(cat "$DEVICE_TOKEN_PATH" 2>/dev/null || true)
    fi

    local payload
    local response
    local next_interval

    payload=$(build_fleet_payload "starting")
    response=$(cloud_post "/api/internal/devices/register" "$payload" "")

    if [ -n "$response" ]; then
        DEVICE_TOKEN=$(printf '%s' "$response" | jq --raw-output '.device_token // empty')
        next_interval=$(printf '%s' "$response" | jq --raw-output '.heartbeat_interval_seconds // empty')

        case "$next_interval" in
            ''|*[!0-9]*)
                ;;
            *)
                HEARTBEAT_INTERVAL="$next_interval"
                ;;
        esac
    fi

    if [ -z "$DEVICE_TOKEN" ] && [ -n "$cached_token" ]; then
        DEVICE_TOKEN="$cached_token"
        bashio::log.warning "Fleet registration unavailable, reusing cached device token."
    fi

    if [ -n "$DEVICE_TOKEN" ]; then
        save_device_token
        bashio::log.info "Fleet reporting active for device uid=${DEVICE_UID}"
    else
        bashio::log.warning "Fleet registration failed, continuing without heartbeat reporting."
        FLEET_ACTIVE=0
    fi
}

send_heartbeat() {
    local state="$1"

    if [ "$FLEET_ACTIVE" -ne 1 ] || [ -z "$DEVICE_TOKEN" ]; then
        return
    fi

    local payload
    local response
    local next_interval

    payload=$(build_fleet_payload "$state")
    response=$(cloud_post "/api/internal/devices/heartbeat" "$payload" "$DEVICE_TOKEN")

    if [ -n "$response" ]; then
        next_interval=$(printf '%s' "$response" | jq --raw-output '.heartbeat_interval_seconds // empty')
        case "$next_interval" in
            ''|*[!0-9]*)
                ;;
            *)
                HEARTBEAT_INTERVAL="$next_interval"
                ;;
        esac
    fi
}

push_ssh_sync_update() {
    local payload="$1"

    if [ "$FLEET_ACTIVE" -ne 1 ] || [ -z "$DEVICE_TOKEN" ]; then
        return
    fi

    cloud_post "/api/internal/devices/ssh-sync" "$payload" "$DEVICE_TOKEN"
}

prepare_sync_payload() {
    local applied_ids_csv="$1"
    local failed_json="$2"
    local revoked_csv="$3"

    jq -n \
        --arg applied "$applied_ids_csv" \
        --arg revoked "$revoked_csv" \
        --argjson failed "$failed_json" \
        '{
            applied_session_ids: ($applied | if . == "" then [] else split(",") | map(tonumber) end),
            failed_sessions: $failed,
            revoked_session_ids: ($revoked | if . == "" then [] else split(",") | map(tonumber) end)
        }'
}

apply_desired_ssh_state() {
    if [ "$FLEET_ACTIVE" -ne 1 ] || [ -z "$DEVICE_TOKEN" ]; then
        return
    fi

    ensure_ssh_storage

    local auth_user
    local auth_keys_path
    local applied_ids_csv
    local initial_payload
    local response
    local new_applied_csv=''
    local new_markers_csv=''
    local revoke_csv=''
    local failed_json='[]'
    local session_id
    local username
    local pubkey
    local fingerprint
    local marker
    local marker_path

    auth_user=$(resolve_ssh_auth_user)
    auth_keys_path=$(get_authorized_keys_path "$auth_user")
    ensure_authorized_keys_parent "$auth_keys_path"

    applied_ids_csv=$(read_applied_ids_csv)
    initial_payload=$(prepare_sync_payload "$applied_ids_csv" '[]' '')
    response=$(push_ssh_sync_update "$initial_payload")

    if [ -z "$response" ]; then
        return
    fi

    local desired_add_count
    local desired_keep_count
    local desired_revoke_count
    desired_add_count=$(printf '%s' "$response" | jq -r '.desired_state.add | length // 0')
    desired_keep_count=$(printf '%s' "$response" | jq -r '.desired_state.keep | length // 0')
    desired_revoke_count=$(printf '%s' "$response" | jq -r '.desired_state.revoke | length // 0')

    if [ "$desired_add_count" -eq 0 ] && [ "$desired_keep_count" -eq 0 ] && [ "$desired_revoke_count" -eq 0 ]; then
        return
    fi

    while IFS=$'\t' read -r session_id username fingerprint pubkey; do
        [ -z "$session_id" ] && continue
        marker=$(marker_name "$session_id" "$fingerprint" "$username")
        marker_path="$SSH_SESSION_MARKERS_PATH/$marker"

        if [ ! -f "$marker_path" ]; then
            if printf '%s\n' "$pubkey" >> "$auth_keys_path"; then
                printf '%s\n' "$pubkey" > "$marker_path"
            else
                failed_json=$(printf '%s' "$failed_json" | jq -c --arg sid "$session_id" --arg err "Unable to append key to authorized_keys" '. + [{session_id: ($sid|tonumber), error: $err}]')
                continue
            fi
        fi

        if [ -z "$new_applied_csv" ]; then
            new_applied_csv="$session_id"
        else
            new_applied_csv="${new_applied_csv},${session_id}"
        fi

        if [ -z "$new_markers_csv" ]; then
            new_markers_csv="$marker"
        else
            new_markers_csv="${new_markers_csv},${marker}"
        fi
    done < <(printf '%s' "$response" | jq -r '.desired_state.add[]? | [.session_id, .ssh_username, (.fingerprint // "none"), .public_key] | @tsv')

    while IFS=$'\t' read -r session_id username fingerprint pubkey; do
        [ -z "$session_id" ] && continue
        marker=$(marker_name "$session_id" "$fingerprint" "$username")
        marker_path="$SSH_SESSION_MARKERS_PATH/$marker"

        if [ ! -f "$marker_path" ]; then
            if printf '%s\n' "$pubkey" >> "$auth_keys_path"; then
                printf '%s\n' "$pubkey" > "$marker_path"
            else
                failed_json=$(printf '%s' "$failed_json" | jq -c --arg sid "$session_id" --arg err "Unable to keep key in authorized_keys" '. + [{session_id: ($sid|tonumber), error: $err}]')
                continue
            fi
        fi

        if [ -z "$new_applied_csv" ]; then
            new_applied_csv="$session_id"
        else
            new_applied_csv="${new_applied_csv},${session_id}"
        fi

        if [ -z "$new_markers_csv" ]; then
            new_markers_csv="$marker"
        else
            new_markers_csv="${new_markers_csv},${marker}"
        fi
    done < <(printf '%s' "$response" | jq -r '.desired_state.keep[]? | [.session_id, .ssh_username, (.fingerprint // "none"), .public_key] | @tsv')

    while IFS=$'\t' read -r session_id username fingerprint; do
        [ -z "$session_id" ] && continue
        marker=$(marker_name "$session_id" "$fingerprint" "$username")
        marker_path="$SSH_SESSION_MARKERS_PATH/$marker"

        if [ -f "$marker_path" ]; then
            pubkey=$(cat "$marker_path" 2>/dev/null || true)
            if [ -n "$pubkey" ]; then
                grep -F -v -- "$pubkey" "$auth_keys_path" > "${auth_keys_path}.tmp" || true
                mv "${auth_keys_path}.tmp" "$auth_keys_path"
            fi
            rm -f "$marker_path"
        fi

        if [ -z "$revoke_csv" ]; then
            revoke_csv="$session_id"
        else
            revoke_csv="${revoke_csv},${session_id}"
        fi
    done < <(printf '%s' "$response" | jq -r '.desired_state.revoke[]? | [.session_id, .ssh_username, (.fingerprint // "none")] | @tsv')

    write_state "$new_applied_csv" "$new_markers_csv"

    local update_payload
    update_payload=$(prepare_sync_payload "$new_applied_csv" "$failed_json" "$revoke_csv")
    push_ssh_sync_update "$update_payload" >/dev/null
}

start_ssh_sync_loop() {
    if [ "$FLEET_ACTIVE" -ne 1 ] || [ -z "$DEVICE_TOKEN" ]; then
        return
    fi

    (
        while true; do
            apply_desired_ssh_state
            sleep "$SSH_SYNC_INTERVAL"
        done
    ) &

    SSH_SYNC_PID=$!
}

start_heartbeat_loop() {
    if [ "$FLEET_ACTIVE" -ne 1 ] || [ -z "$DEVICE_TOKEN" ]; then
        return
    fi

    (
        while true; do
            send_heartbeat "running"
            sleep "$HEARTBEAT_INTERVAL"
        done
    ) &

    HEARTBEAT_PID=$!
}

cleanup() {
    if [ -n "$HEARTBEAT_PID" ]; then
        kill "$HEARTBEAT_PID" 2>/dev/null || true
        wait "$HEARTBEAT_PID" 2>/dev/null || true
    fi

    if [ -n "$SSH_SYNC_PID" ]; then
        kill "$SSH_SYNC_PID" 2>/dev/null || true
        wait "$SSH_SYNC_PID" 2>/dev/null || true
    fi

    if [ "$FLEET_ACTIVE" -eq 1 ] && [ -n "$DEVICE_TOKEN" ]; then
        send_device_log "warn" "device.shutdown" "Cloud Link service stopped" ""
        send_heartbeat "stopped"
    fi
}

trap cleanup EXIT INT TERM

if [ -z "$SUBDOMAIN" ] || [ -z "$ACCESS_TOKEN" ]; then
    bashio::log.error "Subdomain or Access Token is missing"
    bashio::log.error "Please register in Cloud Connect and configure the addon options"
    sleep infinity
fi

if ! [[ "$SUBDOMAIN" =~ ^[a-z0-9-]{3,20}$ ]]; then
    bashio::log.error "Invalid subdomain format. Use 3-20 lowercase letters, digits, or hyphen"
    sleep infinity
fi

if is_valid_port "$SSH_REMOTE_PORT_OPTION"; then
    SSH_REMOTE_PORT="$SSH_REMOTE_PORT_OPTION"
else
    SSH_REMOTE_PORT=0
fi

if is_valid_port "$SSH_LOCAL_PORT_OPTION"; then
    SSH_LOCAL_PORT="$SSH_LOCAL_PORT_OPTION"
else
    SSH_LOCAL_PORT=22
fi

SSH_REMOTE_USER=$(sanitize_remote_user "$SSH_REMOTE_USER_OPTION")
DEVICE_HOSTNAME=$(hostname 2>/dev/null || echo "unknown-device")
DEVICE_UID=$(sanitize_uid "$DEVICE_UID_OPTION")
DEVICE_NAME="$DEVICE_NAME_OPTION"
FLEET_ACTIVE=0

if [ -z "$DEVICE_UID" ]; then
    DEVICE_UID="cloudlink-$(sanitize_uid "$DEVICE_HOSTNAME")-${SUBDOMAIN}"
fi

if [ -z "$DEVICE_NAME" ]; then
    DEVICE_NAME="$DEVICE_HOSTNAME"
fi

if is_truthy "$FLEET_REPORTING_ENABLED" && [ -n "$CLOUD_API_URL" ]; then
    FLEET_ACTIVE=1
fi

bashio::log.info "Initializing Cloud Link for ${SUBDOMAIN}.cloud.apexinfosys.in"

ARCH=$(bashio::info.arch)
FRP_ARCH=""

case $ARCH in
    aarch64) FRP_ARCH="arm64" ;;
    amd64) FRP_ARCH="amd64" ;;
    armhf|armv7) FRP_ARCH="arm" ;;
    i386) FRP_ARCH="386" ;;
    *)
        bashio::log.error "Unsupported architecture: $ARCH"
        sleep infinity
        ;;
esac

if [ -f "$FRPC_BIN" ]; then
    CURRENT_VERSION=$("$FRPC_BIN" -v 2>/dev/null || echo "0")
    if [ "$CURRENT_VERSION" != "$FRP_VERSION" ]; then
        bashio::log.info "FRP version mismatch (${CURRENT_VERSION} vs ${FRP_VERSION}), downloading update"
        rm -f "$FRPC_BIN"
    fi
fi

if [ ! -f "$FRPC_BIN" ]; then
    bashio::log.info "Downloading FRP client for ${FRP_ARCH}"
    mkdir -p "$FRP_DIR"

    wget -qO /tmp/frp.tar.gz "https://github.com/fatedier/frp/releases/download/v${FRP_VERSION}/frp_${FRP_VERSION}_linux_${FRP_ARCH}.tar.gz"
    tar -xzf /tmp/frp.tar.gz -C /tmp/
    mv "/tmp/frp_${FRP_VERSION}_linux_${FRP_ARCH}/frpc" "$FRPC_BIN"
    chmod +x "$FRPC_BIN"

    rm -rf /tmp/frp*
    bashio::log.info "FRP client download complete"
fi

cat << TOML_EOF > "$FRPC_CONFIG"
serverAddr = "${FRP_SERVER_ADDR}"
serverPort = ${FRP_SERVER_PORT}

loginFailExit = false

auth.method = "token"
auth.token = "apexcloud_transport_v1"

[metadatas]
token = "${ACCESS_TOKEN}"

[[proxies]]
name = "ha_web_${SUBDOMAIN}"
type = "http"
localIP = "127.0.0.1"
localPort = 1702
customDomains = ["${SUBDOMAIN}.cloud.apexinfosys.in"]
TOML_EOF

if [ "$SSH_REMOTE_PORT" -gt 0 ]; then
    cat << TOML_SSH_EOF >> "$FRPC_CONFIG"

[[proxies]]
name = "ha_ssh_${SUBDOMAIN}"
type = "tcp"
localIP = "127.0.0.1"
localPort = ${SSH_LOCAL_PORT}
remotePort = ${SSH_REMOTE_PORT}
TOML_SSH_EOF

    bashio::log.info "SSH TCP proxy enabled on ${FRP_SERVER_ADDR}:${SSH_REMOTE_PORT}"
else
    bashio::log.info "SSH TCP proxy is disabled (set ssh_remote_port > 0 to enable)"
fi

if [ "$FLEET_ACTIVE" -eq 1 ]; then
    register_fleet_device
    if [ "$FLEET_ACTIVE" -eq 1 ] && [ -n "$DEVICE_TOKEN" ]; then
        send_device_log "info" "device.start" "Cloud Link service starting" "$(jq -n --arg subdomain "$SUBDOMAIN" --arg frp_addr "$FRP_SERVER_ADDR" --arg frp_port "$FRP_SERVER_PORT" --arg ssh_remote_port "$SSH_REMOTE_PORT" '{subdomain: $subdomain, frp_addr: $frp_addr, frp_port: $frp_port, ssh_remote_port: $ssh_remote_port}')"
        send_heartbeat "starting"
        start_heartbeat_loop
        start_ssh_sync_loop
    fi
else
    bashio::log.info "Fleet reporting disabled"
fi

bashio::log.info "Configuration generated. Connecting to cloud"

set +e
"$FRPC_BIN" -c "$FRPC_CONFIG"
FRPC_EXIT_CODE=$?
set -e

if [ "$FRPC_EXIT_CODE" -eq 0 ]; then
    send_device_log "warn" "frpc.exit" "FRPC process exited cleanly" "$(jq -n --arg code "$FRPC_EXIT_CODE" '{exit_code: $code}')"
else
    send_device_log "error" "frpc.exit" "FRPC process exited with error" "$(jq -n --arg code "$FRPC_EXIT_CODE" '{exit_code: $code}')"
fi

exit "$FRPC_EXIT_CODE"
