# apex-cloud-link/rootfs/etc/services.d/apex-cloud-link/finish
#!/usr/bin/with-contenv bashio
bashio::log.info "Apex Cloud Link service has stopped."
